# config/service/apache2/run
#!/bin/bash

source /etc/apache2/envvars
exec /usr/sbin/apache2ctl -k restart -D FOREGROUND  > /dev/null 2>&1